# rootfs/etc/services.d/kiwix/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Kiwix
# Take care of proper shutdown
# ==============================================================================

bashio::log.info "Kiwix service finished, shutting down addon..."

# Clean up management API if it's running
if [[ -f /tmp/kiwix-manager.pid ]]; then
    MANAGEMENT_PID=$(cat /tmp/kiwix-manager.pid)
    if kill -0 ${MANAGEMENT_PID} 2>/dev/null; then
        bashio::log.info "Stopping management API (PID: ${MANAGEMENT_PID})..."
        kill ${MANAGEMENT_PID} 2>/dev/null || true
    fi
    rm -f /tmp/kiwix-manager.pid
fi

