# rootfs/etc/services.d/kiwix/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Kiwix
# Runs the Kiwix server and management API
# ==============================================================================

bashio::log.info "Starting Kiwix server and management API..."

# Get configuration values
PORT=$(bashio::config 'port')
MANAGEMENT_PORT=$(bashio::config 'management_port')
ZIM_STORAGE_PATH=$(bashio::config 'zim_storage_path')
ENABLE_MANAGEMENT=$(bashio::config 'enable_management')
LOG_LEVEL=$(bashio::config 'log_level')

bashio::log.info "Configuration:"
bashio::log.info "  Kiwix Server Port: ${PORT}"
bashio::log.info "  Management API Port: ${MANAGEMENT_PORT}"
bashio::log.info "  ZIM Storage Path: ${ZIM_STORAGE_PATH}"
bashio::log.info "  Management Enabled: ${ENABLE_MANAGEMENT}"
bashio::log.info "  Log Level: ${LOG_LEVEL}"

# Ensure ZIM storage directory exists and has correct permissions
mkdir -p "${ZIM_STORAGE_PATH}"
chown -R kiwix:kiwix "${ZIM_STORAGE_PATH}"
chmod -R g+w "${ZIM_STORAGE_PATH}"

# Start management API in background if enabled
if bashio::var.true "${ENABLE_MANAGEMENT}"; then
    bashio::log.info "Starting management API on port ${MANAGEMENT_PORT}..."
    export ZIM_STORAGE_PATH="${ZIM_STORAGE_PATH}"
    export MAX_UPLOAD_SIZE=$(bashio::config 'max_upload_size')
    nohup /opt/venv/bin/python /usr/local/bin/kiwix-manager.py \
        --port ${MANAGEMENT_PORT} \
        --storage-path "${ZIM_STORAGE_PATH}" \
        --max-upload-size ${MAX_UPLOAD_SIZE} \
        > /proc/1/fd/1 2>/proc/1/fd/2 &
    MANAGEMENT_PID=$!
    bashio::log.info "Management API started with PID ${MANAGEMENT_PID}"
    
    # Store PID for potential cleanup
    echo ${MANAGEMENT_PID} > /tmp/kiwix-manager.pid
    
    # Wait a moment for management API to start
    sleep 3
    
    # Check if management API is still running
    if ! kill -0 ${MANAGEMENT_PID} 2>/dev/null; then
        bashio::log.error "Management API failed to start!"
    else
        # Verify it's responding
        if curl -f -s http://localhost:${MANAGEMENT_PORT}/api/zim > /dev/null 2>&1; then
            bashio::log.info "Management API is responding on port ${MANAGEMENT_PORT}"
        else
            bashio::log.warning "Management API started but not responding yet. It may need more time."
        fi
    fi
fi

# Start Kiwix server
bashio::log.info "Starting Kiwix server on port ${PORT}..."
bashio::log.info "Kiwix will serve ZIM files from: ${ZIM_STORAGE_PATH}"
bashio::log.info "Using --library and --monitorLibrary for automatic ZIM file detection"

# Check if ingress is enabled and if we need to handle port mismatch
INGRESS_PORT=8111
if [[ "${PORT}" != "${INGRESS_PORT}" ]]; then
    bashio::log.warning "Kiwix configured on port ${PORT}, but ingress expects port ${INGRESS_PORT}"
    bashio::log.warning "Ingress will not work. Use direct access: http://homeassistant-ip:${PORT}"
    bashio::log.warning "Or change port to ${INGRESS_PORT} in configuration for ingress support"
fi

# Check if there are any ZIM files or create empty library
ZIM_COUNT=$(find "${ZIM_STORAGE_PATH}" -maxdepth 1 -name "*.zim" -type f 2>/dev/null | wc -l)
LIBRARY_XML="${ZIM_STORAGE_PATH}/library.xml"

if (( ZIM_COUNT == 0 )); then
    bashio::log.info "No ZIM files found. Creating empty library..."
    # Create an empty library.xml file for kiwix-serve
    cat > "${LIBRARY_XML}" << 'EOF'
<?xml version="1.0" encoding="UTF-8"?>
<library version="2.0" />
EOF
    chown kiwix:kiwix "${LIBRARY_XML}"
    bashio::log.info "Empty library created. Add ZIM files via the management interface."
else
    bashio::log.info "Found ${ZIM_COUNT} ZIM file(s) in storage directory"
fi

# Change to ZIM storage directory
cd "${ZIM_STORAGE_PATH}" || bashio::exit.nok "Could not change to ZIM storage directory"

# Start Kiwix server with library mode and monitoring
# --library: Serve multiple ZIM files
# --monitorLibrary: Automatically detect new/removed ZIM files
# --address: Bind to all interfaces for ingress and direct access
# --port: Listen on configured port
# If library.xml exists, use it; otherwise use current directory
if [[ -f "${LIBRARY_XML}" ]]; then
    exec kiwix-serve --library --monitorLibrary --address 0.0.0.0 --port ${PORT} "${LIBRARY_XML}"
else
    exec kiwix-serve --library --monitorLibrary --address 0.0.0.0 --port ${PORT} .
fi

